# infrastructure/systemd/nexus-backup.service
[Unit]
Description=Nexus encrypted off-host backup
Wants=docker.service
After=docker.service network-online.target

[Service]
Type=oneshot
Environment=ENV_FILE=/etc/nexus/prod.env
Environment=BACKUP_AGE_IDENTITY_FILE=/root/.config/age/nexus-backup-age-identity.txt
WorkingDirectory=/opt/nexus
ExecStart=/bin/sh /opt/nexus/scripts/nexus-backup.sh
UMask=0077
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
