# infrastructure/systemd/nexus-capacity-check.service
[Unit]
Description=Nexus capacity alert check
Wants=docker.service
After=docker.service network-online.target

[Service]
Type=oneshot
Environment=ENV_FILE=/etc/nexus/prod.env
WorkingDirectory=/opt/nexus
ExecStart=/bin/sh /opt/nexus/scripts/check-capacity.sh
UMask=0077
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
