# infrastructure/systemd/nexus-capacity-check.timer
[Unit]
Description=Run Nexus capacity alert check

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
AccuracySec=1min
Persistent=true

[Install]
WantedBy=timers.target
